# scripts/launchd/com.tcground.daily-price-collection.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tcground.daily-price-collection</string>

  <key>WorkingDirectory</key>
  <string>/Users/songjihyeon/VSCodeProjects/tcg-round</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/songjihyeon/.nvm/versions/node/v22.17.0/bin/node</string>
    <string>/Users/songjihyeon/.nvm/versions/node/v22.17.0/lib/node_modules/corepack/dist/pnpm.js</string>
    <string>collect:daily</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/songjihyeon/.nvm/versions/node/v22.17.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/tcground-daily-price-collection.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/tcground-daily-price-collection.err.log</string>
</dict>
</plist>
